# prefect-worker.service
[Unit]
Description=Prefect Worker (serper_tap)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
EnvironmentFile=/opt/serper_tap/.env.production
WorkingDirectory=/opt/serper_tap
ExecStart=/usr/bin/env bash -lc 'cd /opt/serper_tap && source .venv/bin/activate && PREFECT_LOGGING_LEVEL=INFO prefect worker start --pool default-pool'
RestartSec=5
StandardOutput=journal
StandardError=journal
User=svc-prefect
Group=svc-prefect
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
